# roles/rrc-pc/files/pwr-cut-monitor
#!/sbin/openrc-run

name="Pwr-Cut_Monitor"
command="/usr/local/bin/pwr-mgr3.py"
pidfile="/run/pwr-cut-monitor.pid"
command_background=true

