# roles/rrc-pc/files/rc_gnss-share
#!/sbin/openrc-run
# Copyright 2021 Clayton Craft
# SPDX-License-Identifier: GPL-3.0-or-later
description="GNSS location manager"

: ${error_logger="logger -t gnss-share -p daemon.error"}
: ${output_logger="logger -t gnss-share -p daemon.info"}

command=/usr/bin/gnss-share
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"

depend () {
	use logger
}

start_pre() {
    ebegin "Setting GPS receiver to NMEA mode"
    /usr/bin/gpsctl -f -n -s 4800 /dev/ttyUSB_gps
    eend $?
}
